# dreamscene.service
[Unit]
Description=Dream Scene

[Service]
ExecStart=/usr/bin/dreamscene

[Install]
WantedBy=default.target
